# meta-openbmc-machines/meta-openpower/meta-ingrasys/meta-hgx1/recipes-phosphor/obmc-ast-watchdog/files/obmc-ast-watchdog.service
#**************************************************************
#*                                                            *
#*   Copyright (C) Microsoft Corporation. All rights reserved.*
#*                                                            *
#**************************************************************

[Unit]
Description=OpenBMC Hardware Watchdog

[Service]
ExecStartPre=/bin/mkdir -p /run/obmc
ExecStart=@SBINDIR@/obmc-ast-watchdog
Restart=always

[Install]
WantedBy=multi-user.target
