# meta-phosphor/common/recipes-phosphor/fans/obmc-control-fan-algorithm/org.openbmc.control.Fans-Algorithm.service
[Unit]
Description=Phosphor Fan Control

[Service]
Restart=always
ExecStart={sbindir}/fan_algorithm.exe

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
