# meta-openbmc-machines/meta-openpower/meta-ms/meta-sthelens/recipes-phosphor/redfish/obmc-redfish/obmc-redfish.service
[Unit]
Description=OpenBMC console daemon

[Service]
ExecStart=/usr/lib/run_redfish.sh
Restart=always

[Install]
WantedBy=multi-user.target
