# meta-openbmc-machines/meta-openpower/meta-ms/meta-sthelens/recipes-phosphor/redfish/obmc-redfish/obmc-ocslocks-init.service
[Unit]
Description=Init I2C Bus Semaphores

[Service]
Type=oneshot
ExecStart=/usr/sbin/ocslocks-init
Restart=no

[Install]
WantedBy=multi-user.target
