# meta-openbmc-machines/meta-openpower/meta-ms/meta-sthelens/recipes-ipmi/oob-ipmid/oob-ipmid/oob-ipmid.service
[Unit]
Description=Out-Of-Band IPMI Daemon

[Service]
ExecStart=/usr/bin/oob-ipmid
Restart=always

[Install]
WantedBy=multi-user.target
